# src/main/resources/templates/unit.service
# \${unit.fileName}.service generated by @project.artifactId@ v@project.version@

[Unit]
Description=\${unit.description}
After=\${unit.after}

[Service]
EnvironmentFile=\${environment.fileName}
WorkingDirectory=\${unit.workingDirectory}
Type=\${unit.type}
ExecStart=\${unit.execStart}
Restart=\${unit.restart}
User=\${unit.user}

[Install]
WantedBy=${unit.wantedBy}
